# mpd-ratings-load.path
[Path]
PathChanged=%h/.config/mpd/ratings_sync
TriggerLimitBurst=0

[Install]
WantedBy=default.target